# templates/etcd.service.erb
[Unit]
Description=etcd key-value store
Documentation=https://github.com/etcd-io/etcd
After=network.target

[Service]
User=<%= @user %>
Group=<%= @group %>
Type=notify
ExecStart=<%= @bin_dir %>/etcd --config-file <%= @config_path %>
Restart=always
RestartSec=10s
LimitNOFILE=<%= @max_open_files %>

[Install]
WantedBy=multi-user.target
